# deploy/systemd/paisareality-prices.timer
# Runs the daily price update at 07:00 in the server timezone (IST).
# Persistent=true catches up the run if the server was off at the scheduled time.
[Unit]
Description=Run Paisa Reality price update daily

[Timer]
OnCalendar=*-*-* 07:00:00
Persistent=true

[Install]
WantedBy=timers.target
